# a4/tarecho@.service
[Unit]
Description=Tarecho service

[Service]
ExecStart=/usr/local/bin/compile_troff.py
StandardInput=socket
